# src/Test_RPI/bootprog.sh
#!/bin/sh -e
#
### BEGIN INIT INFO
# Provides:          StartUp
# Required-Start:    $all
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: description du programme
### END INIT INFO

cd /home/user/ //C program path
./Hello //Program name
